# docker/7.4/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autorestart=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=true
